# emuld.service
[Unit]
Description=Daemon for emulator
After=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/emuld

[Install]
WantedBy=emulator.target

