# emuld.service
[Unit]
Description=Common control daemon for emulator
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/emuld
