# emuld.service
[Unit]
Description=emuld
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/emuld
